# roles/haproxy/templates/le-renew-haproxy.timer.j2
[Unit]
Description=Let's encrypt renew HAProxy

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=multi-user.target
